# php-cli/supervisor.ini
[program:crond]
command=/usr/sbin/crond -f
stdout_logfile=/dev/sdtout
stderr_logfile=/dev/sdterr

[group:core]
programs=crond
